# crontab
# Runs process manager jobs
* * * * * php /var/www/html/bin/console process-manager:maintenance